# ansible/roles/tailscale/files/tailscale-cert-renew.timer
[Unit]
Description=Daily Tailscale TLS certificate renewal

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
