# config/supervisord-base.conf
[supervisord]
nodaemon=true
logfile=/dev/stdout
logfile_maxbytes=0

[program:xvnc]
command=Xvnc :%(ENV_DISPLAY_NUM)s -geometry %(ENV_DISPLAY_RESOLUTION)s -depth 24 -rfbport %(ENV_VNC_PORT)s -rfbauth /tmp/vncpasswd -localhost no -AlwaysShared
autorestart=true
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:cursor]
command=bash -c "sleep 2 && DISPLAY=:%(ENV_DISPLAY_NUM)s xsetroot -cursor_name left_ptr"
autorestart=false
startsecs=0
startretries=3
priority=15
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:scrcpy]
command=scrcpy --no-audio --max-fps %(ENV_SCRCPY_MAX_FPS)s -b %(ENV_SCRCPY_BITRATE)s --window-width %(ENV_DISPLAY_W)s --window-height %(ENV_DISPLAY_H)s --window-borderless --window-x 0 --window-y 0
environment=DISPLAY=":%(ENV_DISPLAY_NUM)s"
autorestart=true
startsecs=5
startretries=999
priority=20
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:gnirehtet]
command=gnirehtet run -d %(ENV_GNIREHTET_DNS)s
autorestart=true
startsecs=5
startretries=999
priority=20
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:audio-server]
command=/usr/local/bin/otacon-server
autorestart=true
startsecs=3
startretries=999
priority=30
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
