# config/supervisord-bluetooth.conf

[program:dbus]
command=bash -c 'rm -f /run/dbus/pid /run/dbus/system_bus_socket && dbus-daemon --system --nofork'
autorestart=true
priority=1
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:bluetoothd]
command=bluetoothd -n -E --noplugin=sap
autorestart=true
priority=2
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:bluealsa]
; HFP HF mode: Pi acts as a Bluetooth headset (hands-free unit).
; The phone connects as AG and routes mic audio over SCO to BlueALSA.
; Sleep 3s so bluetoothd is ready to register the HFP HF profile.
command=bash -c 'sleep 3 && bluealsa --profile=hfp-hf --profile=a2dp-sink'
autorestart=true
priority=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:bt-agent]
; Sleep 5s so BlueALSA+bluetoothd are ready before agent registers
command=bash -c 'sleep 5 && python3 /opt/bluetooth-agent.py'
autorestart=true
priority=4
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
